# templates/etc/init/keystone.conf
description "OpenStack Keystone service"

start on runlevel [2345]
stop on runlevel [016]

setuid {{keystone.user}}
respawn
exec /opt/stack/keystone/bin/keystone-all

